# modules/minecraft/templates/minecraft_init.erb
#!/bin/bash
# /etc/init.d/minecraft

  ### BEGIN INIT INFO
  # Provides:   minecraft
  # Required-Start: $local_fs $remote_fs
  # Required-Stop:  $local_fs $remote_fs
  # Should-Start:   $network
  # Should-Stop:    $network
  # Default-Start:  2 3 4 5
  # Default-Stop:   0 1 6
  # Short-Description:    Minecraft server
  # Description:    Starts the minecraft server
  ### END INIT INFO

#Settings
USERNAME="<%= @user %>"
MCPATH="<%= @homedir %>"
MARK2=/usr/local/bin/mark2

ME=`whoami`
as_user() {
  if [ "$ME" == "$USERNAME" ] ; then
    bash -c "$1"
  else
    su - $USERNAME -c "$1"
  fi
}

mc_start() {
  as_user "cd $MCPATH ; timeout 20s $MARK2 start"
}

mc_stop() {
  as_user "cd $MCPATH ; $MARK2 send save-all ; sleep 10 ; $MARK2 stop"
}

#Start-Stop here
case "$1" in
  start)
    mc_start
    ;;
  stop)
    mc_stop
    ;;
  restart)
    mc_stop
    mc_start
    ;;
  status)
    if ps ax | grep -v grep | grep spigot > /dev/null
    then
      echo "Minecraft is running."
    else
      echo "Minecraft is not running."
      exit 1
    fi
    ;;

  *)
  echo "Usage: /etc/init.d/minecraft {start|stop|status|restart}"
  exit 1
  ;;
esac

exit 0
